# TraderEngine.CLI/.assets/container/root/etc/crontabs/root
 2 * * * * dotnet /app/TraderEngine.CLI.dll -automations